# deploy/systemd/xml-import.path
[Unit]
Description=Watch XML import directory

[Path]
PathExistsGlob=/opt/xml_data_importer/*.xml

[Install]
WantedBy=multi-user.target
